# system/cron.d/container-healthcheck
*/5 * * * * root /opt/containers/healthcheck.sh
