# system/cron.d/finance-scraper
# Scrape Aegon Retiready pension value daily at 7am
0 7 * * * root cd /opt/containers/finance-scraper && docker-compose run --rm finance-scraper >> /var/log/finance-scraper.log 2>&1
